# images/zookeeper/resources/supervisor/zookeeper.conf
[program:zookeeper]
command=python /opt/zookeeper-3.4.6/pod/pod.py
